# contrib/init/luascoind.conf
description "LUASCOIN Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env LUASCOIND_BIN="/usr/bin/luascoind"
env LUASCOIND_USER="luascoincore"
env LUASCOIND_GROUP="luascoincore"
env LUASCOIND_PIDDIR="/var/run/luascoind"
# upstart can't handle variables constructed with other variables
env LUASCOIND_PIDFILE="/var/run/luascoind/luascoind.pid"
env LUASCOIND_CONFIGFILE="/etc/luascoincore/luascoin.conf"
env LUASCOIND_DATADIR="/var/lib/luascoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # luascoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$LUASCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run luascoind."
        echo "The setting must appear in $LUASCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"LUASCOIN Core Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$LUASCOIND_PIDDIR"
    chmod 0755 "$LUASCOIND_PIDDIR"
    chown $LUASCOIND_USER:$LUASCOIND_GROUP "$LUASCOIND_PIDDIR"
    chown $LUASCOIND_USER:$LUASCOIND_GROUP "$LUASCOIND_CONFIGFILE"
    chmod 0660 "$LUASCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$LUASCOIND_PIDFILE" \
    --chuid $LUASCOIND_USER:$LUASCOIND_GROUP \
    --exec "$LUASCOIND_BIN" \
    -- \
    -pid="$LUASCOIND_PIDFILE" \
    -conf="$LUASCOIND_CONFIGFILE" \
    -datadir="$LUASCOIND_DATADIR" \
    -disablewallet \
    -daemon

